# files/system/usr/lib/systemd/system/unbound.service.d/secureblue.conf
# SPDX-FileCopyrightText: Copyright (c) 2007, NLnet Labs. All rights reserved.
# SPDX-FileCopyrightText: Copyright 2025-2026 The Secureblue Authors
#
# SPDX-License-Identifier: BSD-3-Clause
#
# Modified from: https://github.com/NLnetLabs/unbound/blob/65be1d0ec30986362996f81549c6c630616d08f7/contrib/unbound.service.in

# Use systemd sockets instead of CAP_NET_BIND_SERVICE.
[Unit]
Requires=unbound.socket unbound-control.socket
After=unbound.socket unbound-control.socket

[Service]
# Usually starts as root:root, setuid/setgid -> unbound:unbound.
User=unbound
Group=unbound

# Upstream was CAP_NET_BIND_SERVICE CAP_SETGID CAP_SETUID CAP_SYS_CHROOT
# CAP_SYS_RESOURCE CAP_NET_RAW.
CapabilityBoundingSet=
ProtectSystem=strict
ReadWritePaths=/run/unbound /etc/unbound /var/lib/unbound
InaccessiblePaths=/run/dbus/ /run/user/
RuntimeDirectory=unbound
RuntimeDirectoryMode=0775
RuntimeDirectoryPreserve=yes
ConfigurationDirectory=unbound
StateDirectory=unbound
SystemCallFilter=@system-service

# Taken from unbound/contrib upstream.
# systemd merges both SystemCallFilters, with the first (i.e. the allowlist above)
# dictating the initial action.
SystemCallFilter=~@clock @cpu-emulation @debug @keyring @module mount @obsolete @resources
MemoryDenyWriteExecute=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=false
ProtectProc=invisible
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
RestrictRealtime=true
SystemCallArchitectures=native
RestrictNamespaces=yes
LockPersonality=yes
RestrictSUIDSGID=yes
