# files/system/usr/lib/systemd/system/dnsconfd.service.d/secureblue.conf
# SPDX-FileCopyrightText: Copyright 2025-2026 The Secureblue Authors
#
# SPDX-License-Identifier: Apache-2.0

[Service]
# SELinux policy allows net_bind_service net_raw net_admin, supposedly for
# unbound-control.
CapabilityBoundingSet=
ProtectSystem=strict
ReadWritePaths=/etc/resolv.conf /run/dnsconfd
ProtectHome=yes
PrivateTmp=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectProc=invisible
LockPersonality=yes
RestrictNamespaces=yes
PrivateDevices=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
RestrictAddressFamilies=AF_INET AF_INET6 AF_NETLINK AF_UNIX
PrivateMounts=yes
IPAddressDeny=any
IPAddressAllow=localhost
RestrictRealtime=yes
